# home/dot_config/systemd/user/tmux-restore.service
[Unit]
Description=Restore saved tmux 'main' session (after the server is up)
After=tmux.service
Wants=tmux.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Login shell: gives restore.sh the full PATH/env it needs (tmux, plugins).
ExecStart=/bin/bash -lc '/home/ze-tank/.config/tmux/boot-restore.sh'

[Install]
WantedBy=default.target
